# service/apache/run
#!/bin/sh
exec /usr/sbin/apache2ctl -D FOREGROUND 2>&1
